# Docker/supervisord.conf
[supervisord]
nodaemon=true
user=root
loglevel=info

# Asterisk (foreground obligatoire avec -f)
[program:asterisk]
command=/usr/sbin/asterisk -f -vvv
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0

# EarlyWarning
# sleep : laisse Asterisk démarrer l'AMI avant la 1re connexion de l'appli
[program:earlywarning]
command=/bin/sh -c "sleep 4 && exec java $JAVA_OPTS -jar /earlywarning/EarlyWarning.jar"
directory=/earlywarning
autostart=true
autorestart=true
priority=20
stdout_logfile=/dev/stdout
stdout_logfile_maxbytes=0
stderr_logfile=/dev/stderr
stderr_logfile_maxbytes=0